# sd.ini
[program:sd]
command=/usr/bin/php /data/sd/bin/start_swoole_server.php start
autostart=true
startsecs=3
autorestart=true
log_stderr=true
log_stdout=true
log_stderr=true
logfile=/var/log/supervisor/sd.log
